# systemd/dspilot-concessions-daily.service
[Unit]
Description=DSPilot — Daily concessions + investigations pipeline
After=network-online.target dspilot-report-daily.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/DSPilot

ExecStart=/bin/bash /root/DSPilot/scripts/run-concessions-pipeline.sh

TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
