# systemd/dspilot-investigations-weekly.service
[Unit]
Description=DSPilot — Weekly investigations scrape from S3 + ingestion
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/DSPilot
EnvironmentFile=/root/.secrets/dspilot.env
Environment=NEXT_PUBLIC_CONVEX_URL=https://sincere-rhinoceros-718.convex.cloud

ExecStart=/bin/bash -c '\
  cd /root/DSPilot && \
  set -a && source /root/.secrets/dspilot.env && set +a && \
  python3 scraper/amazon_investigations_sync.py \
    --company-id 55c476fc-c3d4-4577-9910-586110f01405 \
    --station-filter DIF1 \
    --invoke-ingest \
    --station-code FR-PSUA-DIF1 \
    --organization-id org_37Yb7MlFJHFs5h7K28zFYoZ4fUY'

TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
